# site-template/cdbweb/cdbweb.service.in
[Unit]
Description=CDB Web Service
Documentation=@CDB_WEB_DOCURL@
After=network.target 

[Service]

User=@USERID@
Group=@GROUPID@

Environment="PYTHONPATH=@PYTHONPATH@"
Environment="CDB_ROOT_DIR=@CDB_ROOT_DIR@"
Environment="PYTHON_EGG_CACHE=@PYTHON_EGG_CACHE@"

## We don't need to add --daemon option, becuase systemd will make this application as a daemon.
## We also don't need to have a pid file, but just in case we would like to use `cdb-web-service.sh`.
##
ExecStart=@PYTHON_CMD@ @CDB_DAEMON_CMD@ --pid-file @CDB_WEB_PID_FILE@ --config-file @CDB_WEB_CONFIG_FILE@

Restart=always

[Install]
WantedBy=multi-user.target
#